# deploy/telegram-managed-bot.service.example
# Example systemd unit for a Telegram bot managed by CreaBotManager.
# Copy to /etc/systemd/system/telegram-my-bot.service and adjust paths.
#
# TimeoutStopSec=10 avoids the default 90s wait on stop/restart when the process
# does not exit quickly on SIGTERM. For instant shutdown, handle SIGTERM in the bot.

[Unit]
Description=My Telegram Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=crearec
WorkingDirectory=/home/crearec/my-bot
EnvironmentFile=/home/crearec/my-bot/.env
ExecStart=/usr/bin/env node dist/index.js
Restart=always
RestartSec=5
TimeoutStopSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=telegram-my-bot

[Install]
WantedBy=multi-user.target
